# packaging/pkgmgr-installer-signal-agent.socket
[Socket]
ListenStream=/run/pkgmgr/agent/%U
SocketMode=0777
DirectoryMode=0777

[Install]
WantedBy=sockets.target
